# systemd/meowtrics.service
[Unit]
Description=meowtrics tray pet daemon
Documentation=https://github.com/ra-yavuz/meowtrics
PartOf=graphical-session.target
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/meowtrics daemon
Restart=on-failure
RestartSec=5
# Security hardening. The daemon only reads /proc and /sys and renders a tray icon.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=true
RestrictNamespaces=true
RestrictRealtime=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

[Install]
WantedBy=graphical-session.target
