# systemd/meowtrics.service
[Unit]
Description=meowtrics tray pet daemon
Documentation=https://github.com/ra-yavuz/meowtrics
PartOf=graphical-session.target
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/meowtrics daemon
Restart=on-failure
RestartSec=5

# This is a user-scope unit. Most of systemd's sandboxing primitives
# (ProtectSystem, ProtectHome, PrivateTmp, MemoryDenyWriteExecute,
# SystemCallFilter) require CAP_SYS_ADMIN to apply and fail with
# status=218/CAPABILITIES when run via the user manager. We keep only the
# ones that are safe in user mode.
NoNewPrivileges=true
RestrictRealtime=true
LockPersonality=true

[Install]
WantedBy=graphical-session.target
